# IRCNotificationBot.service
[Unit]
Description=Simple IRC notifier
After=syslog.target

[Service]
Type=simple
Restart=always
ExecStart=/opt/IRCNotificationBot/watchdog.py
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
